# a2i-devops/stage/ansible/roles/nifi-toolkit/templates/ca_server.service.tpl
Description= This service is used for controlling ca-server

Wants=network.target
After=syslog.target network-online.target

[Service]
Type=simple
ExecStart=/usr/local/nifi-toolkit/bin/enable-server.sh
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target
